# .config/systemd/user/keep.target
[Unit]
Description=save
[Install]
Alias=default.target
