# .config/systemd/user/monkey.target
[Unit]
Description=monkey
[Install]
Alias=default.target
