# .config/systemd/user/scrutor.target
[Unit]
Description=search
[Install]
Alias=default.target
